# provisioning/roles/nfs/files/nfs-downloads.automount
[Unit]
Description=mount /nfs/downloads

[Automount]
Where=/nfs/downloads
TimeoutIdleSec=300

[Install]
WantedBy=multi-user.target